# deploy/ansible/templates/kanban.service.j2
[Unit]
Description=Kanban Runtime (Phoung)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory={{ kanban_repo_dir }}
ExecStart=/usr/bin/node dist/cli.js --host 127.0.0.1 --port {{ kanban_port }} --no-open
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile={{ kanban_env_file }}

[Install]
WantedBy=multi-user.target
